# go/k8s/install/k8s-hard-way/ansible-demo/demo-1/roles/etcd/templates/etcd.service.j2


[Unit]
Description=etcd
After=network.target

[Service]
ExecStart=/usr/bin/etcd --name {{ ansible_default_ipv4.address }} --addr {{ ansible_default_ipv4.address }}:{{ etcd.port }} --data-dir /var/lib/etcd --bind-addr 0.0.0.0:{{ etcd.port }}
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
